# ais_forwarder.service
[Unit]
Description=AIS Forwarder Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ais_forwarder
Restart=on-failure
RestartSec=30
StartLimitInterval=300
StartLimitBurst=5
User=david
Group=david
Environment=DISPLAY=:0
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ais_forwarder
KillMode=control-group
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
